# docker/with-nginx/supervisor.conf
[supervisord]
nodaemon=true
[program:gawk]
command=/app/bin/start.sh
autostart=true
autorestart=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[program:nginx]
command=/app/docker/nginx/start-nginx.sh
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
logfile_backups=0
